# etc/rc.local
#!/bin/bash
#
# /etc/rc.local: Local multi-user startup script.
#
mount tmpfs /tmp -t tmpfs -o size=8192m
